# pkg/pftbld.rc
#!/bin/ksh

daemon="/usr/local/sbin/pftbld"

. /etc/rc.d/rc.subr

rc_cmd $1
